# server/scripts/headlenss.service.template
[Unit]
Description=headlenss server (tmux + ASR + Web UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{SERVER_DIR}}
EnvironmentFile=-{{SERVER_DIR}}/.env
ExecStart={{NPM}} start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Make sure tmux sessions created by the server share the user's default tmux server
Environment=PATH={{PATH}}

[Install]
WantedBy=default.target
